# templates/timelapse-sync.timer.jinja2
[Unit]
Description=Timelapse upload every 5 minutes

[Timer]
Unit=timelapse.service
OnCalendar=*:0/5
RandomizedDelaySec=60

[Install]
WantedBy=timers.target
