# templates/timelapse-sync.service.jinja2
[Unit]
Description=Timelapse picture upload

[Service]
Type=oneshot
ExecStart=/opt/timelapse/timelapse-sync.sh 

[Install]
WantedBy=multi-user.target
